# sftpgo/sftpgo.files/B_DEBIAN/DEBIAN/prerm
#!/bin/sh

systemctl stop sftpgo.service
systemctl disable sftpgo.service
deluser --quiet sftpgo
