# kanboard/kanboard.files/B_DEBIAN/etc/cron.daily/kanboard
#!/bin/sh
cd /opt/kanboard && ./cli cronjob >/dev/null 2>&1
