# powerdns-admin/powerdns-admin.files/B_DEBIAN/etc/systemd/system/powerdns-admin.service.d/override.conf
[Service]
Environment="FLASK_CONF=../configs/production.py"