# authelia/authelia.files/B_DEBIAN/lib/systemd/system/authelia@.service
[Unit]
Description=Authelia authentication and authorization server
After=multi-user.target

[Service]
User=authelia
Group=authelia
ExecStart=/opt/authelia/bin/authelia --config /etc/authelia-%i/configuration.yml
SyslogIdentifier=authelia-%i
Restart=always

[Install]
WantedBy=multi-user.target
