# dmarc-smtptls-reports/dmarc-smtptls-reports.files/B_DEBIAN/etc/cron.daily/dmarc-smtptls-reports
#!/bin/bash
cd /opt/dmarc-smtptls-reports
php bin/console app:getreportsfrommailbox
