# jitsi-excalidraw-backend/jitsi-excalidraw-backend.files/B_DEBIAN/DEBIAN/prerm
#!/bin/sh

systemctl stop jitsi-excalidraw-backend.service
systemctl disable jitsi-excalidraw-backend.service
rm -rf /opt/excalidraw-backend/dist
deluser --quiet excalidraw
