# gvm/gsa.files/B_DEBIAN/DEBIAN/postinst
systemctl enable gsa
systemctl start gsa
