# snipe-it/snipe-it.files/B_DEBIAN/etc/cron.d/snipe-it
* * * * * www-data /usr/bin/php /opt/snipe-it/artisan schedule:run >> /dev/null 2>&1
