# illuminanced/illuminanced.files/B_DEBIAN/DEBIAN/postinst
#!/bin/sh

set -e

systemctl enable illuminanced.service
