# powerdns/pdns.files/B_DEBIAN/DEBIAN/prerm
#!/bin/sh

systemctl stop pdns.service
systemctl disable pdns.service
deluser --quiet pdns
