# cryptpad/cryptpad.files/B_DEBIAN/DEBIAN/prerm
#!/bin/sh

set -e
systemctl stop cryptpad || true
systemctl disable cryptpad || true
userdel --force cryptpad > /dev/null || true
groupdel cryptpad > /dev/null || true
