# openwebui/openwebui.files/B_DEBIAN/DEBIAN/prerm
#!/bin/sh

systemctl stop openwebui.service
systemctl disable openwebui.service
deluser --quiet openwebui
