# nginx/nginx.files/B_DEBIAN/DEBIAN/prerm
#!/bin/sh
systemctl stop nginx
systemctl disable nginx
