# iventoy/iventoy.files/B_DEBIAN/DEBIAN/prerm
#!/bin/sh

systemctl stop iventoy
systemctl disable iventoy
