# gvm/openvas.files/B_DEBIAN-BUSTER/etc/cron.daily/openvas
#!/bin/sh
su gvm -c /usr/bin/openvas-feed-update
