# corefreq/corefreq.files/B_DEBIAN/DEBIAN/postinst
#!/bin/sh
 
if test "$1" = "configure"; then
  systemctl enable corefreqd.service
  systemctl start corefreqd.service
else
  systemctl restart corefreqd.service
fi
