# corefreq/corefreq.files/B_DEBIAN/DEBIAN/prerm
#!/bin/sh
 
systemctl stop corefreqd.service
systemctl disable corefreqd.service
