# corefreq/corefreqk-dkms.files/B_DEBIAN/DEBIAN/postinst
dkms add -m corefreqk -v 1.98.2
dkms build -m corefreqk -v 1.98.2 && dkms install -m corefreqk -v 1.98.2 || true
if test "" = "configure"; then
  systemctl daemon-reload
  systemctl enable corefreqd.service
  systemctl start corefreqd.service
else
  systemctl daemon-reload
  systemctl restart corefreqd.service
fi
