# corefreq/corefreqk-dkms.files/B_DEBIAN/DEBIAN/postinst
systemctl stop corefreqd.service || true
rmmod corefreqk || true
dkms add -m corefreqk -v 2.0.0
dkms build -m corefreqk -v 2.0.0 && dkms install -m corefreqk -v 2.0.0
modprobe corefreqk
systemctl enable corefreqd.service
systemctl daemon-reload
systemctl start corefreqd.service
