# corefreq/corefreqk-dkms.files/B_DEBIAN/DEBIAN/postinst
systemctl stop corefreqd.service || true
rmmod corefreqk || true
dkms add -m corefreqk -v 1.98.5
dkms build -m corefreqk -v 1.98.5 && dkms install -m corefreqk -v 1.98.5
modprobe corefreqk
systemctl enable corefreqd.service
systemctl daemon-reload
systemctl start corefreqd.service
