# corefreq/corefreqk-dkms.files/B_DEBIAN/DEBIAN/prerm
systemctl stop corefreqd.service
systemctl disable corefreqd.service
rmmod corefreqk || true
dkms remove -m corefreqk -v 2.0.0 --all || true
