# hollama/hollama.files/B_DEBIAN/DEBIAN/prerm
#!/bin/sh

systemctl stop hollama.service
systemctl disable hollama.service
deluser --quiet hollama
