# ollama/ollama.files/B_DEBIAN/DEBIAN/prerm
#!/bin/sh

systemctl stop ollama.service
systemctl disable ollama.service
deluser --quiet ollama
